# service/smirc.service
[Unit]
Description=System monitor IRC bot

[Service]
Type=simple
ExecStart=/usr/bin/smirc --bot

[Install]
WantedBy=multi-user.target

